# csye6225.service
[Unit]
Description=CSYE6225 Web Application
After=network.target

[Service]
User=csye6225
Group=csye6225
WorkingDirectory=/opt/csye6225/webapp
ExecStart=/usr/bin/node app.js
Restart=always
EnvironmentFile=/opt/csye6225/webapp/.env


[Install]
WantedBy=multi-user.target

# Reload systemd and enable/start service
sudo systemctl daemon-reload
sudo systemctl enable csye6225.service
sudo systemctl start csye6225.service